# packaging/debian/prerm
#!/bin/sh
set -e

case "${1:-}" in
	remove | deconfigure)
		guard=/usr/lib/steward-node/release/scripts/node-removal-guard.sh
		[ -x "$guard" ] || guard=/opt/steward/current/integration/scripts/node-removal-guard.sh
		[ -x "$guard" ] || guard=/usr/local/libexec/steward/node-removal-guard
		"$guard"
		systemctl stop steward-gateway.service >/dev/null 2>&1 || true
		systemctl stop steward.service >/dev/null 2>&1 || true
		systemctl stop steward-executor.service >/dev/null 2>&1 || true
		systemctl disable steward-gateway.service steward.service steward-executor.service >/dev/null 2>&1 || true
		;;
esac

exit 0
