# packaging/debian/prerm
#!/bin/sh
set -e

case "${1:-}" in
	remove | deconfigure)
		systemctl disable --now steward.service steward-executor.service steward-gateway.service >/dev/null 2>&1 || true
		;;
esac

exit 0
