# deploy/systemd/steward-control.service
[Unit]
Description=Steward fleet control plane
Documentation=https://github.com/hardrails/steward
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=steward-control
Group=steward-control
UMask=0077
EnvironmentFile=/etc/steward-control/control.env
ExecStart=/usr/local/bin/steward-control -addr=${STEWARD_CONTROL_ADDR} -state-dir=${STEWARD_CONTROL_STATE_DIR} -auth-key-file=${STEWARD_CONTROL_AUTH_KEY_FILE} -witness-private-key-file=${STEWARD_CONTROL_WITNESS_PRIVATE_KEY_FILE} -witness-public-key-file=${STEWARD_CONTROL_WITNESS_PUBLIC_KEY_FILE} -tls-cert-file=${STEWARD_CONTROL_TLS_CERT_FILE} -tls-key-file=${STEWARD_CONTROL_TLS_KEY_FILE} -enable-metrics=${STEWARD_CONTROL_ENABLE_METRICS} -node-stale-after=${STEWARD_CONTROL_NODE_STALE_AFTER} -evidence-stale-after=${STEWARD_CONTROL_EVIDENCE_STALE_AFTER} -command-overdue-after=${STEWARD_CONTROL_COMMAND_OVERDUE_AFTER} -capacity-warning-percent=${STEWARD_CONTROL_CAPACITY_WARNING_PERCENT}
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
StateDirectory=steward-control
StateDirectoryMode=0700
TasksMax=256
LimitNOFILE=4096
MemoryMax=1G
MemorySwapMax=0

NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=
PrivateDevices=yes
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
LockPersonality=yes
MemoryDenyWriteExecute=yes
RemoveIPC=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
